# debian/cgroupfs-mount.upstart
description "mount available cgroup filesystems"
author "Serge Hallyn <serge.hallyn@canonical.com>"

start on mounted MOUNTPOINT=/sys/fs/cgroup
stop on unmounted-remote-filesystems

pre-start script
	test -x /usr/bin/cgroupfs-mount || { stop; exit 0; }
	/usr/bin/cgroupfs-mount
end script

post-stop script
	test -x /usr/bin/cgroupfs-umount || exit 0
	/usr/bin/cgroupfs-umount
end script
